# vowels/rootfs/etc/services.d/vowels/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Vowels (Automated Electrical Energy System)..."

# Run our .NET Native AOT daemon
exec /app/vowels
